# scripts/picap.service
[Unit]
Description=PiCap camera OCR service (HTTP + BLE)
After=network-online.target bluetooth.service
Wants=network-online.target bluetooth.service

[Service]
Type=simple
User=woodworn
Group=woodworn
WorkingDirectory=/home/woodworn/PiCap
Environment=PATH=/home/woodworn/PiCap/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PICAP_USER=woodworn

# Wait for hci0 + enable LE (as root). Timed; never blocks forever.
ExecStartPre=-+/home/woodworn/PiCap/scripts/picap-bt-prep.sh

# Starts pairing agent, then PiCap serve.
ExecStart=/home/woodworn/PiCap/scripts/picap-systemd-start.sh

Restart=on-failure
RestartSec=5

# Prep may wait up to ~40s for hci0 after power-on.
TimeoutStartSec=90
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL

SupplementaryGroups=bluetooth video

[Install]
WantedBy=multi-user.target
